# ndntdump-capture-start.service
[Unit]
Description=Start Capturing NDN Traffic
After=docker.service network-online.target

[Service]
ExecStart=/usr/local/bin/ndntdump-capture-start.sh

[Install]
WantedBy=multi-user.target