# ndntdump-scp.service
[Unit]
Description=Copy Traces To Tntech Node
After=ndntdump-capture-stop.service

[Service]
ExecStart=/usr/local/bin/ndntdump-scp.sh